# piped.service
[Unit]
Description=Tobacco Pipe Webscraper

[Service]
ExecStart=/usr/local/bin/piped

[Install]
WantedBy=multi-user.target
